# scripts/cosmos_video_downloader.service
[Unit]
Description=Cosmos Football Video Downloader
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cosmos-football-event-detection
ExecStart=/usr/bin/python3 /home/ubuntu/cosmos-football-event-detection/scripts/video_download_scheduler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
